# cloude-bot.service.example
# Servicio systemd de USUARIO para cloude-bot.
# Instalación:
#   mkdir -p ~/.config/systemd/user
#   cp cloude-bot.service.example ~/.config/systemd/user/cloude-bot.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now cloude-bot.service
#   loginctl enable-linger $USER   # para que siga tras cerrar sesión
#
# Logs:   journalctl --user -u cloude-bot -f

[Unit]
Description=cloude-bot - Claude Code por Telegram
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/valle/proyectos/cloude-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/valle/proyectos/cloude-bot/.venv/bin/python /home/valle/proyectos/cloude-bot/src/telegram_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
